# etc/service/lightKnot/run
#!/bin/sh
chpst -ug3d2 /home/seb/repos/lightKnot/lightKnot.js
