# etc/service/http/run
#!/bin/sh
chpst -uroot /opt/http/server.js
